# deploy/systemd/trading-engine.service
[Unit]
Description=C++ Trading Engine — Multi-venue SOR (Binance/Kraken/OKX/Coinbase)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/trading
EnvironmentFile=/etc/trading/env
ExecStart=/opt/trading/build/bin/trading_engine --mode live --venues BINANCE,KRAKEN,OKX,COINBASE --symbol BTCUSDT --loop-interval-ms 500
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=trading-engine

[Install]
WantedBy=multi-user.target
